# scripts/cloudquiz.service




[Unit]
Description=Curio Cloud Quiz — FastAPI/Gunicorn Application
Documentation=https://github.com/rohitrk/curio
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify

# ── Runtime user ─────────────────────────────────────────────
# Never run as root
User=ubuntu
Group=ubuntu

# ── Working directory ─────────────────────────────────────────
# MUST be project root (THE2/ equivalent) so backend.app.* imports resolve
WorkingDirectory=/home/ubuntu/curio

# ── Environment ───────────────────────────────────────────────
# Load .env file (pydantic-settings also reads this, but systemd loads it too
# so gunicorn subprocess has all vars available)
EnvironmentFile=/home/ubuntu/curio/backend/.env

# PYTHONPATH: lets Python find backend.app.* from the project root
Environment=PYTHONPATH=/home/ubuntu/curio

# ── Gunicorn command ──────────────────────────────────────────
# Worker count guidance:
#   t3.micro  (1 vCPU,  1 GB RAM) → -w 2
#   t3.small  (2 vCPU,  2 GB RAM) → -w 4
#   t3.medium (2 vCPU,  4 GB RAM) → -w 4
#   t3.large  (2 vCPU,  8 GB RAM) → -w 5
#
# Module: backend.app.main:app
# Matches development command: uvicorn backend.app.main:app --reload
ExecStart=/home/ubuntu/curio/venv/bin/gunicorn backend.app.main:app \
    --worker-class  uvicorn.workers.UvicornWorker \
    --workers       2 \
    --bind          127.0.0.1:8000 \
    --timeout       120 \
    --keep-alive    5 \
    --max-requests  1000 \
    --max-requests-jitter 50 \
    --log-level     info \
    --access-logfile  /var/log/curio/access.log \
    --error-logfile   /var/log/curio/error.log \
    --capture-output

# ── Graceful reload (used by: systemctl reload curio) ─────────
# Sends SIGHUP to Gunicorn master — workers drain in-flight
# requests before restarting (true zero-downtime)
ExecReload=/bin/kill -s HUP $MAINPID

# ── Restart policy ────────────────────────────────────────────
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=3

# ── Resource limits ───────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096

# ── Security hardening ────────────────────────────────────────
# Prevent privilege escalation
NoNewPrivileges=true
# Restrict what the process can write to (except /home/ubuntu/curio and /var/log/curio)
ProtectSystem=full
# Keep /home readable so the app can access its own files
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/curio/static /var/log/curio /tmp
# Prevent access to /dev devices
PrivateDevices=true
# Isolate temp files per service
PrivateTmp=true

[Install]
WantedBy=multi-user.target
