# scripts/cloudquiz.service
# /etc/systemd/system/cloudquiz.service
#
# Manages the Cloud Quiz FastAPI/Gunicorn process.
# Commands:
#   sudo systemctl start   cloudquiz
#   sudo systemctl stop    cloudquiz
#   sudo systemctl restart cloudquiz
#   sudo systemctl status  cloudquiz
#   sudo journalctl -u cloudquiz -f          # live logs
# ─────────────────────────────────────────────────────────────

#[Unit]
#Description=Cloud Quiz FastAPI Application
#After=network.target postgresql.service
#Wants=postgresql.service

#[Service]
#Type=notify

# Run as the ubuntu user (not root)
#User=ubuntu
#Group=ubuntu

# Working directory — adjust if you deploy elsewhere
#WorkingDirectory=/home/ubuntu/cloud_quiz/backend

# Load .env automatically
#EnvironmentFile=/home/ubuntu/cloud_quiz/backend/.env

# Gunicorn with Uvicorn workers
#   -w 2  →  2 × (2 × CPU cores + 1) is a common starting point
#            For a t3.micro (1 vCPU) use -w 2
#            For a t3.small/medium (2 vCPU) use -w 4
#ExecStart=/home/ubuntu/cloud_quiz/venv/bin/gunicorn app.main:app \
#    --worker-class uvicorn.workers.UvicornWorker \
#    --workers 2 \
#    --bind 127.0.0.1:8000 \
#    --timeout 120 \
#    --keep-alive 5 \
#    --log-level info \
#    --access-logfile /var/log/cloudquiz/access.log \
#    --error-logfile  /var/log/cloudquiz/error.log

# Graceful reload on SIGHUP (used by systemctl reload)
#ExecReload=/bin/kill -s HUP $MAINPID

# Restart policy
#Restart=on-failure
#RestartSec=5s

# File descriptor limits
#LimitNOFILE=65536

#[Install]
#WantedBy=multi-user.target
















# /etc/systemd/system/curio.service
#
# Manages the Curio FastAPI application via Gunicorn + Uvicorn workers.
#
# Commands:
#   sudo systemctl start   curio
#   sudo systemctl stop    curio
#   sudo systemctl restart curio
#   sudo systemctl reload  curio       ← graceful reload (zero-downtime)
#   sudo systemctl status  curio
#   sudo journalctl -u curio -f        ← live logs

[Unit]
Description=Curio Cloud Quiz — FastAPI Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify

# Run as the ubuntu user (never root)
User=ubuntu
Group=ubuntu

# !! CRITICAL: WorkingDirectory must be the project ROOT (THE2/ equivalent)
# so that Python resolves  backend.app.*  imports correctly.
WorkingDirectory=/home/ubuntu/curio

# Load .env automatically into the process environment
EnvironmentFile=/home/ubuntu/curio/backend/.env

# Add the project root to PYTHONPATH so  backend.app.*  imports work
Environment=PYTHONPATH=/home/ubuntu/curio

# Gunicorn with Uvicorn workers
#   -w 2 → good for t3.micro (1 vCPU)
#   -w 4 → better for t3.small/medium (2 vCPU)
#
# Module path: backend.app.main:app
# (matches:  uvicorn backend.app.main:app --reload  used in development)
ExecStart=/home/ubuntu/curio/venv/bin/gunicorn backend.app.main:app \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 2 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --keep-alive 5 \
    --log-level info \
    --access-logfile /var/log/curio/access.log \
    --error-logfile  /var/log/curio/error.log

# Graceful reload on SIGHUP — used by  systemctl reload curio
ExecReload=/bin/kill -s HUP $MAINPID

# Restart on crash, but not on clean stop
Restart=on-failure
RestartSec=5s

# Allow large numbers of open files (needed for many concurrent connections)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
